# .docker/etc/crontabs/root
* * * * * php /www/artisan schedule:run >> /dev/null 2>&1